# roles/base/files/systemd/cpulimited.slice
[Unit]
Description=Slice used to run daily clamscan. 30% CPU
Before=slices.target

[Slice]
CPUQuota=30%

[Install]
WantedBy=multi-user.target
