# roles/base/files/systemd/daily-clamscan.timer
[Unit]
Description=Run daily Clamscan at 14:00

[Timer]
OnCalendar=*-*-* 14:00
Persistent=true
Unit=daily-clamscan.service

[Install]
WantedBy=timers.target
