# roles/base/files/systemd/freshclam-updatedb.service
[Unit]
Description=Weekly Freshclam updateDB
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/freshclam --quiet
