# roles/base/files/systemd/logfile-reduction.timer
[Unit]
Description=Run logfile reduction after boot

[Timer]
OnBootSec=20
Unit=logfile-reduction.service

[Install]
WantedBy=timers.target
